# thriller-bark/backups/the-sea-dumps.timer
[Unit]
Description=Run nightly DB dumps at 03:00 UTC

[Timer]
OnCalendar=*-*-* 03:00:00 UTC
Persistent=true

[Install]
WantedBy=timers.target
